# root/etc/services.d/php7-fpm/run
#!/bin/sh

cd /srv/www/

if [ -z "$(ls -A /srv/www/wordpress)" ]; then
    tar -xJf wordpress.tar.xz
fi

cp wp-config.php wordpress/

exec php-fpm7 --nodaemonize --fpm-config /etc/php7/php-fpm.conf
