# root/etc/services.d/nginx/run
#!/bin/sh

exec s6-setuidgid nginx \
    nginx -g 'daemon off;' -c /etc/nginx/nginx.conf